# resources/docker/services.d/nginx/run
#!/bin/sh
nginx -g "daemon off;"
echo "nginx running"